# clash/clash.service
[Unit]
Description=Clash Board

[Service]
ExecStart=/usr/bin/clash -d /etc/clash

[Install]
WantedBy=multi-user.target

